# roles/shadowsocks/templates/shadowsocks.service.j2
[Unit]
Description=Shadowsocks

[Service]
TimeoutStartSec=0
ExecStart={{ ssserver.stdout }} -c /etc/shadowsocks.json

[Install]
WantedBy=multi-user.target
